# deploy/livekit-agent.service
[Unit]
Description=LiveKit outbound voice agent (Nora)
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/live-kit-tuts
EnvironmentFile=/home/ubuntu/live-kit-tuts/.env
ExecStart=/home/ubuntu/.local/bin/uv run python agent.py start
Restart=on-failure
RestartSec=3
StandardOutput=append:/var/log/livekit-agent.log
StandardError=append:/var/log/livekit-agent.log

[Install]
WantedBy=multi-user.target
